# proof-of-concept/coreos/pcr-rebind-boot.service
[Unit]
After=cryptsetup.target
[Service]
Type=oneshot
ExecStart=/usr/local/bin/rebindDiskOnBoot.sh
[Install]
WantedBy=multi-user.target
